# systemd/garmin-sync.service
[Unit]
Description=Garmin Playwright Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=florent
WorkingDirectory=/home/florent/garmin-playwright-sync
ExecStart=/home/florent/garmin-playwright-sync/run.sh
TimeoutStartSec=120
